# server/systemd_service/aicamera-frontend.service
[Unit]
Description=AI Camera Frontend Build Service
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/devuser/aicamera/server/frontend
Environment=PATH=/usr/bin:/usr/local/bin:/home/devuser/.npm-global/bin
ExecStart=/usr/bin/npm run build
ExecReload=/usr/bin/npm run build
# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aicamera-frontend

[Install]
WantedBy=default.target