# server/systemd_service/aicamera-backend.service
[Unit]
Description=AI Camera Server Backend (NestJS)
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/home/devuser/aicamera/server
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/bin:/usr/local/bin:/home/devuser/.npm-global/bin
ExecStart=/usr/bin/node /home/devuser/aicamera/server/dist/src/main.js
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3
# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aicamera-backend

[Install]
WantedBy=default.target