# server/systemd_service/backend-api.service
[Unit]
Description=AI Camera Backend API
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=devuser
WorkingDirectory=/home/devuser/aicamera/server/backend-api
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target


