# kiosk-browser-headless.service
[Unit]
Description=AI Camera Browser (Headless Mode - No Display Required)
After=network.target aicamera_lpr.service nginx.service
Wants=network.target
Requires=network.target aicamera_lpr.service nginx.service

[Service]
Type=simple
User=camuser
Group=camuser
WorkingDirectory=/home/camuser/aicamera

# ไม่ต้องการ DISPLAY และ XAUTHORITY ใน headless mode
# Environment=DISPLAY=:0
# Environment=XAUTHORITY=/home/camuser/.Xauthority

# Wait for the web service to be ready before starting browser
# Check both nginx and the Flask app are responding
ExecStartPre=/bin/bash -c 'echo "Waiting for nginx to be ready..."'
ExecStartPre=/bin/bash -c 'timeout 30 bash -c "until systemctl is-active --quiet nginx; do sleep 2; done"'
ExecStartPre=/bin/bash -c 'echo "Waiting for aicamera_lpr service to be ready..."'
ExecStartPre=/bin/bash -c 'timeout 60 bash -c "until systemctl is-active --quiet aicamera_lpr.service; do sleep 3; done"'
ExecStartPre=/bin/bash -c 'echo "Waiting for web service to respond..."'
ExecStartPre=/bin/bash -c 'timeout 120 bash -c "until curl -s -f http://localhost/ >/dev/null 2>&1; do sleep 5; done"'
ExecStartPre=/bin/bash -c 'echo "Web service is ready!"'
ExecStart=/bin/bash -c 'echo "Starting browser in headless mode..." && sleep 5 && chromium-browser --headless --disable-gpu --no-sandbox --disable-dev-shm-usage --disable-web-security --disable-features=VizDisplayCompositor --no-first-run --no-default-browser-check --disable-background-timer-throttling --disable-backgrounding-occluded-windows --disable-renderer-backgrounding --disable-features=TranslateUI --disable-ipc-flooding-protection --autoplay-policy=no-user-gesture-required --disable-session-crashed-bubble --disable-infobars --disable-notifications --new-window --allow-running-insecure-content --user-data-dir=/tmp/chromium-headless --remote-debugging-port=9222 http://localhost/'

# Restart policy - restart if browser exits
Restart=always
RestartSec=30

# Kill the browser process on service stop
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-camera-browser-headless

[Install]
WantedBy=multi-user.target
