# edge/systemd_service/kiosk-browser-normal.service
[Unit]
Description=AI Camera Browser (Normal Mode - No Fullscreen)
After=network.target aicamera_lpr.service nginx.service
Wants=network.target
Requires=network.target aicamera_lpr.service nginx.service

[Service]
Type=simple
User=camuser
Group=camuser
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/camuser/.Xauthority
WorkingDirectory=/home/camuser/aicamera

# Wait for the web service to be ready before starting browser
# Check both nginx and the Flask app are responding
ExecStartPre=/bin/bash -c 'echo "Waiting for nginx to be ready..."'
ExecStartPre=/bin/bash -c 'timeout 30 bash -c "until systemctl is-active --quiet nginx; do sleep 2; done"'
ExecStartPre=/bin/bash -c 'echo "Waiting for aicamera_lpr service to be ready..."'
ExecStartPre=/bin/bash -c 'timeout 60 bash -c "until systemctl is-active --quiet aicamera_lpr.service; do sleep 3; done"'
ExecStartPre=/bin/bash -c 'echo "Waiting for web service to respond..."'
ExecStartPre=/bin/bash -c 'timeout 120 bash -c "until curl -s -f http://localhost/ >/dev/null 2>&1; do sleep 5; done"'
ExecStartPre=/bin/bash -c 'echo "Web service is ready!"'
ExecStart=/bin/bash -c 'echo "Starting browser in normal mode..." && sleep 5 && chromium-browser --user-data-dir=/tmp/chromium-normal --app=http://localhost/ --no-first-run --no-default-browser-check --disable-infobars --disable-notifications'

# Restart policy - restart if browser exits
Restart=always
RestartSec=30

# Kill the browser process on service stop
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-camera-browser

[Install]
WantedBy=multi-user.target